# scripts/crontab
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin
* * * * *  DISPLAY=:0 /home/arduino/DaemonFMDDatalogger.sh